# deploy/claude-bridge.service
[Unit]
Description=Claude Bridge - IM to Claude Code
After=network.target

[Service]
Type=simple
User=xiang
WorkingDirectory=/home/xiang
ExecStart=/home/xiang/.local/bin/uv run --project /home/xiang/code/claude-bridge python -m claude_bridge
Restart=always
RestartSec=5
Environment=HOME=/home/xiang
EnvironmentFile=-/home/xiang/.claude-bridge/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
