# deploy/uclaw.service
[Unit]
Description=uclaw - Turn Claude Code into OpenClaw
After=network.target

[Service]
Type=simple
User=xiang
WorkingDirectory=/home/xiang
ExecStart=/home/xiang/.local/bin/uv run --project /home/xiang/code/uclaw python -m uclaw
Restart=always
RestartSec=5
Environment=HOME=/home/xiang
EnvironmentFile=-/home/xiang/.uclaw/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
